# jarvis.service
[Unit]
Description=JARVIS Caps Lock PTT Voice Assistant
After=graphical-session.target

[Service]
ExecStart=/home/kisa/jarvis-env/bin/python3 -u /home/kisa/jarvis/main_final.py
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus

[Install]
WantedBy=default.target
